# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:nginx]
command=/usr/local/openresty/nginx/sbin/nginx
